# configs/musix/root-image/etc/rc.local
#!/bin/bash
#
# /etc/rc.local: Local multi-user startup script.
#
# set passwd
sed -i 's/^arch:[^:]*/arch:$6$O3ydLHFN$15KdwteKFNRgoe06rieoiiCl9cY\/JWjLSBeoExo2ef3YKFMcXGEXuRBNQJAO.ax3QbuKbidnXED0HqcZ4FetG\./' /etc/shadow

#su arch -c mpd
su arch -c jackstart
alsactl restore
