# configs/musix/root-image/etc/rc.local
#!/bin/bash
#
# /etc/rc.local: Local multi-user startup script.
#
su radio -c mpd
su radio -c jackstart
alsactl restore
